# teleghost.service.example
[Unit]
Description=TeleGhost — Telegram ↔ Mattermost transparent bridge
After=network.target

[Service]
Type=simple
User=YOUR_USER
WorkingDirectory=/path/to/teleghost
ExecStart=/path/to/teleghost/.venv/bin/python3 -m teleghost
Environment=PYTHONPATH=/path/to/teleghost/src
Environment=TELEGHOST_CONFIG=/path/to/teleghost/config.yaml
Restart=always
RestartSec=5
MemoryMax=512M

[Install]
WantedBy=multi-user.target
